# src/mysql/build/modules/build/files/etc/supervisor/conf.d/mysql_stdout.conf
[program:mysql_stdout]
command = /usr/bin/tail -F /var/log/mysql/stdout.log
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
