# src/mysqld/build/modules/build/files/etc/supervisor/conf.d/mysqld_stderr.conf
[program:mysqld_stderr]
command = /usr/bin/tail -F /var/log/mysql/stderr.log
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
