# src/build/modules/mysqld/files/etc/supervisor/conf.d/mysqld.conf
[program:mysqld]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
