# build/etc/supervisor/conf.d/mysql.conf
[program:mysql]
command=mysqld_safe
autostart=true
autorestart=true