# build/etc/supervisor/conf.d/mysql.conf
[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true